# deploy/systemd/ai-video-daily-report-check.timer
[Unit]
Description=Check AI Video Daily Report approval replies every 5 minutes

[Timer]
OnCalendar=*:0/5
AccuracySec=30s
Persistent=true
Unit=ai-video-daily-report-check.service

[Install]
WantedBy=timers.target
