# deploy/systemd/dingtalk-order-poller.service
[Unit]
Description=DingTalk supplier order poller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/vipthink/dingtalk-order-bot
EnvironmentFile=/etc/vipthink/dingtalk-order-bot.env
Environment=HOME=/root
ExecStart=/opt/vipthink/dingtalk-order-bot/.venv/bin/python scripts/poll_dingtalk_group_orders.py --interval 30 --limit 100 --reply
Restart=always
RestartSec=10
StandardOutput=append:/var/log/dingtalk-order-poller.log
StandardError=append:/var/log/dingtalk-order-poller.err.log

[Install]
WantedBy=multi-user.target
