# launchd/com.vipthink.global-ig-scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vipthink.global-ig-scheduler</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/yangyi/.vipthink-auto-publish/launchd/run_vipthink_global_ig_scheduler.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yangyi/.vipthink-auto-publish</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>18</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/yangyi/.vipthink-auto-publish/logs/com.vipthink.global-ig-scheduler.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/yangyi/.vipthink-auto-publish/logs/com.vipthink.global-ig-scheduler.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
